# dot_config/systemd/user/notify-shutdown.timer
[Unit]
Description=Timer for shutdown notification

[Timer]
OnCalendar=*-*-* 22:50:00
Persistent=false

[Install]
WantedBy=timers.target
