# dot_config/systemd/user/notify-shutdown.service
[Unit]
Description=Popup shutdown notification

[Service]
Type=oneshot
ExecStart=/usr/bin/notify-send "System will shutdown at 22:58"
