# dot_config/systemd/user/sync-ventoy.service
[Unit]
Description=Sync Documents and Music to Ventoy USB when mounted

[Service]
Type=oneshot
ExecStart=/home/nazze/sh/ventoy_sync.sh
Restart=no
