# dot_config/systemd/user/sync-ventoy.path
[Unit]
Description=Watch Ventoy mount point for sync

[Path]
PathExists=/run/media/nazze/Ventoy
Unit=sync-ventoy.service

[Install]
WantedBy=default.target
