# dot_config/systemd/user/udiskie.service
[Unit]
Description=udiskie automounter
After=graphical-session.target

[Service]
ExecStart=/usr/bin/udiskie
Restart=on-failure
Type=simple

[Install]
WantedBy=default.target
