# dot_config/systemd/root/plan-poweroff.service
[Unit]
Description=Scheduled shutdown at 22:58

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl poweroff
