# dot_config/systemd/root/plan-poweroff.timer
[Unit]
Description=Timer to poweroff at 22:55 daily

[Timer]
OnCalendar=*-*-* 22:55:00

[Install]
WantedBy=timers.target
